# marklogic-data-hub-central/src/rpm/scripts/post-install.sh
#!/bin/sh
echo "Enabling service: hub-central"
systemctl enable hub-central